# systemd/blitz-bot-manager.service
[Unit]
Description=Blitz Trading Bot Manager
After=network.target blitz-web.service
Wants=network.target
Requires=blitz-web.service

[Service]
Type=simple
User=blitz
Group=blitz
WorkingDirectory=/opt/blitz-test-server
Environment=DATABASE_URL=sqlite:////opt/blitz-test-server/instance/production.db
ExecStart=/opt/blitz-test-server/venv/bin/python /opt/blitz-test-server/bot_manager.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/blitz-test-server/instance
ReadWritePaths=/opt/blitz-test-server/logs

[Install]
WantedBy=multi-user.target